# src/autoinit/cron.d/autoinit
@reboot root /bin/bash /root/autoinit/init_script
